# etc/udev/rules.d/95-hotplug-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", HOTPLUG=="1", RUN+="/home/eli/.config/hotplug-monitor.sh"
